# ssd-temp.service
[Unit]
Description=ssd-temp

[Service]
Type=simple
ExecStart=-/opt/ssd-temp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
